# deploy/trumptrade.service
[Unit]
Description=TrumpTrade Bot
Documentation=https://github.com/YOUR_USERNAME/trump_trade
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trumptrade
Group=trumptrade
WorkingDirectory=/home/trumptrade/trumptrade
EnvironmentFile=/home/trumptrade/trumptrade/.env
ExecStart=/home/trumptrade/trumptrade/.venv/bin/trumptrade
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trumptrade
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
